# addon-AqaraPost/rootfs/etc/s6-overlay/s6-rc.d/init-nodered/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Node-RED
# Configures Node-RED before running
# ==============================================================================
declare port
declare renew_flow=false
declare flow_workdir=""
declare aqara_default_server=""
declare aqara_default_appid=""
declare aqara_default_appkey=""

escape_sed_value() {
	printf '%s' "$1" | sed -e 's/[\\/&]/\\&/g'
}

set_aqara_defaults() {
	case "${1^^}" in
		CN)
			aqara_default_server="aiot-rpc.aqara.cn"
			aqara_default_appid="94549908487478b220992a70"
			aqara_default_appkey="Jddz01kIORDYrBzqGYgpUXKBnIHfW8E3"
			;;
		EU|AF)
			aqara_default_server="rpc-ger.aqara.com"
			aqara_default_appid="7be1984f0556276133336839"
			aqara_default_appkey="Jddz01kIORDYrBzqGYgpUXKBnIHfW8E3"
			;;
		RU)
			aqara_default_server="rpc-ru.aqara.com"
			aqara_default_appid="94549908487478b220992a70"
			aqara_default_appkey="euGhPe2rcmxwculATNj45eEtnd50zp0I"
			;;
		KR|JP)
			aqara_default_server="rpc-kr.aqara.com"
			aqara_default_appid="94549908487478b220992a70"
			aqara_default_appkey="euGhPe2rcmxwculATNj45eEtnd50zp0I"
			;;
		AU|ME)
			aqara_default_server="rpc-au.aqara.com"
			aqara_default_appid="94549908487478b220992a70"
			aqara_default_appkey="Jddz01kIORDYrBzqGYgpUXKBnIHfW8E3"
			;;
		USA|US|HMT|OTHER|*)
			aqara_default_server="aiot-rpc-usa.aqara.com"
			aqara_default_appid="94549908487478b220992a70"
			aqara_default_appkey="Jddz01kIORDYrBzqGYgpUXKBnIHfW8E3"
			;;
	esac
}

apply_flow_values() {
	local flow_file=$1
	local lumi_value=$2

	sed -i \
		-e "s|XXXXXXAREAXXXXXXXXXXX|$(escape_sed_value "$AREA")|g" \
		-e "s|XXXXXXTOKENXXXXXXXXXXX|$(escape_sed_value "$TOKEN")|g" \
		-e "s|XXXXXXAPPIDXXXXXXXXXXX|$(escape_sed_value "$APPID")|g" \
		-e "s|XXXXXXAPPKEYXXXXXXXXXXX|$(escape_sed_value "$APPKEY")|g" \
		-e "s|XXXXXXUSERIDXXXXXXXXXXX|$(escape_sed_value "$USERID")|g" \
		-e "s|lumi1.XXXXXXXXXXXX|$(escape_sed_value "$lumi_value")|g" \
		-e "s|rpc-ger.aqara.com|$(escape_sed_value "$SERVER")|g" \
		-e "s|it-IT|$(escape_sed_value "$TIMEZONE")|g" \
		"$flow_file"
}

if bashio::config.true 'renew_flow'; then
	bashio::log.info "renew_flow true"
	renew_flow=true
else
	bashio::log.info "renew_flow false"
fi

	#Install node-forge
	bashio::log.info "......."
	bashio::log.info "Install node-forge ..."
    cd /opt/node_modules && npm install node-forge

	#Install cryptojs
	bashio::log.info "......."
	bashio::log.info "Install cryptojs ..."
    cd /opt/node_modules && npm install cryptojs

if bashio::fs.file_exists '/config/flows.json' && [ "$renew_flow" != "true" ]; then
	bashio::log.info "......."
    bashio::log.info "flows.json already exist, the file will not be updated from config"
	bashio::log.info "......."
else
	USERNAME=$(bashio::config 'usernameAqara') || USERNAME=
	PASSWORD=$(bashio::config 'passwordAqara') || PASSWORD=

	if [ "$USERNAME" = "" ] || [ "$USERNAME" = "preconfigured" ] || [ "$PASSWORD" = "" ] || [ "$PASSWORD" = "preconfigured" ] || [ "$PASSWORD" = "[]" ]; then
		if bashio::fs.file_exists '/config/flows.json'; then
			bashio::log.warning "renew_flow requested but Aqara credentials are missing or masked; keeping existing flows.json"
			renew_flow=false
		else
			bashio::exit.nok "Cannot create AqaraPOST flows: usernameAqara/passwordAqara are missing or masked"
		fi
	fi

fi

if [ "$renew_flow" = "true" ] || ! bashio::fs.file_exists '/config/flows.json'; then
	#Install python for token generator aqara
	bashio::log.info "......."
	bashio::log.info "Install python3/pycryptodome for token generator Aqara ..."
	apk add python3 && apk add py3-pip && pip install requests && pip install pycryptodome && pip install requests-toolbelt  && bashio::log.info "Install python3/pycryptodome complete"|| bashio::log.info "install python3 error"
	
	#AqaraPost
	bashio::log.info "......."
	flow_workdir=$(mktemp -d) \
		|| bashio::exit.nok "Failed to create temporary AqaraPost flow folder"
	bashio::log.info "Prepare AqaraPost flow in temporary folder"
	bashio::log.info "......."
	
	bashio::log.info "Download json Aqara from GitHub"
	touch /etc/node-red/flows.json || bashio::log.info "...."
	bashio::log.info "......."

	FLOW_SOURCE_BASE_URL=$(bashio::config 'flow_source_base_url') || bashio::log.info "error flow_source_base_url ......."
	if [ "$FLOW_SOURCE_BASE_URL" = "" ]; then
	FLOW_SOURCE_BASE_URL="https://raw.githubusercontent.com/sdavides/AqaraPOST-Homeassistant/main"
	fi
	FLOW_SOURCE_BASE_URL="${FLOW_SOURCE_BASE_URL%/}"
	bashio::log.info "Download AqaraPost flow files from $FLOW_SOURCE_BASE_URL"

	FLOW_URL1="$FLOW_SOURCE_BASE_URL/Aqara_G3_nodered.json"
	FLOW_NAME1="Aqara_G3_nodered"

	FLOW_URL2="$FLOW_SOURCE_BASE_URL/Aqara_FP2_nodered.json"
	FLOW_NAME2="Aqara_FP2_nodered"

	FLOW_URL3="$FLOW_SOURCE_BASE_URL/Aqara_G3_nodered_2device.json"
	FLOW_NAME3="Aqara_G3_nodered_2device"

	LUMI1=$(bashio::config 'lumi1Aqara') || bashio::log.info "error lumi1Aqara ......."
	LUMI2=$(bashio::config 'lumi1Aqara2') || bashio::log.info "error lumi1Aqara2 ......."
	LUMI3=$(bashio::config 'lumi1Aqara3') || bashio::log.info "error lumi1Aqara3 ......."

	TIMEZONE=$(bashio::config 'timezoneAqara') || bashio::log.info "error timezoneAqara ......."

	AREA=$(bashio::config 'serverAqara') || bashio::log.info "error serverAqara ......."
	set_aqara_defaults "$AREA"

    if [ "$LUMI1" = "" ]; then	
	bashio::log.warning "G3 lumi1 not set"
	else
	LUMI1=$(echo "$LUMI1" | tr '[:upper:]' '[:lower:]' | sed -E 's/-//g; /^lumi1\./!s/^/lumi1./')
	wget -q "$FLOW_URL1" --output-document="$flow_workdir/$FLOW_NAME1.json" && bashio::log.info "Download AqaraPost $FLOW_NAME1.json complete..." || bashio::exit.nok "Error download AqaraPost $FLOW_NAME1 flow json"
	fi
	
    if [ "$LUMI2" = "" ]; then
	bashio::log.warning "FP2 lumi1 not set"
	else		
	LUMI2=$(echo "$LUMI2" | tr '[:upper:]' '[:lower:]' | sed -E 's/-//g; /^lumi1\./!s/^/lumi1./')
	wget -q "$FLOW_URL2" --output-document="$flow_workdir/$FLOW_NAME2.json" && bashio::log.info "Download AqaraPost $FLOW_NAME2.json complete..." || bashio::exit.nok "Error download AqaraPost $FLOW_NAME2 flow json"
	fi	
	
	if [ "$LUMI3" = "" ]; then	
	bashio::log.warning "G3_2 lumi1 not set"
	else
	LUMI3=$(echo "$LUMI3" | tr '[:upper:]' '[:lower:]' | sed -E 's/-//g; /^lumi1\./!s/^/lumi1./')
	wget -q "$FLOW_URL3" --output-document="$flow_workdir/$FLOW_NAME3.json" && bashio::log.info "Download AqaraPost $FLOW_NAME3.json complete..." || bashio::exit.nok "Error download AqaraPost $FLOW_NAME3 flow json"
	fi
	
	ls "$flow_workdir"/Aqara_*.json || bashio::exit.nok "No AqaraPost flow json was downloaded"
		
	bashio::log.info "Starting AqaraPost-Homeassistant tokenGenerator ..."
	chmod +x /usr/local/bin/AqaraPOST-tokenGenerator.py || bashio::log.info "tokenGenerator script already executable"
	( echo "$USERNAME" ; echo "$PASSWORD"; echo "$(bashio::config 'serverAqara')" ) | /usr/local/bin/AqaraPOST-tokenGenerator.py > /tmp/AqaraPOST-tokenGenerator.json || bashio::log.warning "tokenGenerator command exited with an error"
	cat /tmp/AqaraPOST-tokenGenerator.json | grep Username || bashio::log.info "Username Token error......."
	cat /tmp/AqaraPOST-tokenGenerator.json | grep Area || bashio::log.info "Area Token error......."

    if grep -q "Token:" /tmp/AqaraPOST-tokenGenerator.json; then
    TOKEN=$(grep "Token:" /tmp/AqaraPOST-tokenGenerator.json | sed 's/Token://')
    else
    bashio::log.warning "error Token..."
	TOKEN=
    fi

    if grep -q "AppID:" /tmp/AqaraPOST-tokenGenerator.json; then
    APPID=$(grep "AppID:" /tmp/AqaraPOST-tokenGenerator.json | sed 's/AppID://')
    else
    bashio::log.warning "error AppID..."
	APPID=
    fi

    if grep -q "AppKey:" /tmp/AqaraPOST-tokenGenerator.json; then
    APPKEY=$(grep "AppKey:" /tmp/AqaraPOST-tokenGenerator.json | sed 's/AppKey://')
    else
    bashio::log.warning "error AppKey..."
	APPKEY=
    fi

    if grep -q "UserID:" /tmp/AqaraPOST-tokenGenerator.json; then
    USERID=$(grep "UserID:" /tmp/AqaraPOST-tokenGenerator.json | sed 's/UserID://')
    else
    bashio::log.warning "error UserID..."
	USERID=
    fi

    if grep -q "Server:" /tmp/AqaraPOST-tokenGenerator.json; then
    SERVER=$(grep "Server:" /tmp/AqaraPOST-tokenGenerator.json | sed 's/Server://')
	cat /tmp/AqaraPOST-tokenGenerator.json | grep Server
    else
    bashio::log.warning "error Server..."
	SERVER=
    fi

	if [ "$APPID" = "" ]; then
		APPID=$aqara_default_appid
		bashio::log.warning "Using static AppID for $AREA"
	fi

	if [ "$APPKEY" = "" ]; then
		APPKEY=$aqara_default_appkey
		bashio::log.warning "Using static AppKey for $AREA"
	fi

	if [ "$SERVER" = "" ]; then
		SERVER=$aqara_default_server
		bashio::log.warning "Using static server for $AREA"
	fi

	if [ "$TOKEN" = "" ] || [ "$APPID" = "" ] || [ "$APPKEY" = "" ] || [ "$USERID" = "" ] || [ "$SERVER" = "" ]; then
		if bashio::fs.file_exists '/config/flows.json'; then
			bashio::log.warning "tokenGenerator did not return all required values; keeping existing flows.json"
			rm -rf "$flow_workdir" || bashio::log.info "......."
			rm -rf /tmp/AqaraPOST-tokenGenerator.json || bashio::log.info "......."
			renew_flow=false
		else
			bashio::log.warning "tokenGenerator did not return all required values; generating Node-RED flow with placeholder Aqara credentials so the add-on can still start"
			if [ "$TOKEN" = "" ]; then
				TOKEN="MISSING_AQARA_TOKEN"
			fi
			if [ "$USERID" = "" ]; then
				USERID="MISSING_AQARA_USERID"
			fi
			if [ "$APPID" = "" ]; then
				APPID="MISSING_AQARA_APPID"
			fi
			if [ "$APPKEY" = "" ]; then
				APPKEY="MISSING_AQARA_APPKEY"
			fi
			if [ "$SERVER" = "" ]; then
				SERVER="rpc-ger.aqara.com"
			fi
		fi
	fi

	if [ "$renew_flow" = "true" ] || ! bashio::fs.file_exists '/config/flows.json'; then
	
	if  bashio::fs.file_exists "$flow_workdir/$FLOW_NAME1.json"; then
	apply_flow_values "$flow_workdir/$FLOW_NAME1.json" "$LUMI1" && bashio::log.info "tokenGenerator Insert value $FLOW_NAME1.json complete." || bashio::exit.nok "Error AqaraPost-Homeassistant tokenGenerator"
	fi
		
	if  bashio::fs.file_exists "$flow_workdir/$FLOW_NAME2.json"; then
	apply_flow_values "$flow_workdir/$FLOW_NAME2.json" "$LUMI2" && bashio::log.info "tokenGenerator Insert value $FLOW_NAME2.json complete." || bashio::exit.nok "Error AqaraPost-Homeassistant tokenGenerator"
	fi
	
	if  bashio::fs.file_exists "$flow_workdir/$FLOW_NAME3.json"; then
	apply_flow_values "$flow_workdir/$FLOW_NAME3.json" "$LUMI3" && bashio::log.info "tokenGenerator Insert value $FLOW_NAME3.json complete." || bashio::exit.nok "Error AqaraPost-Homeassistant tokenGenerator"
	fi
		
	jq -s 'add' "$flow_workdir"/Aqara_*json > /tmp/AqaraPOST-flows.json && mv /tmp/AqaraPOST-flows.json /config/flows.json && bashio::log.info "flows.json created" || bashio::exit.nok "Error creating flows.json"
	rm -rf "$flow_workdir" || bashio::log.info "......."
	rm -rf /tmp/AqaraPOST-tokenGenerator.json || bashio::log.info "......."

	fi
	
	bashio::log.info "......."
	bashio::log.info "End AqaraPost-Homeassistant flow" 
	bashio::log.info "......."

fi

	bashio::log.info "Starting Node-Red ..." 
	
# Migrate add-on data from the Home Assistant config folder,
# to the add-on configuration folder.
if ! bashio::fs.file_exists '/config/settings.js' \
    && bashio::fs.file_exists '/homeassistant/node-red/settings.js'; then
    shopt -s dotglob
    mv /homeassistant/node-red/* /config/ \
        || bashio::exit.nok "Failed to migrate Node-RED configuration"
fi

# Ensure configuration exists
if ! bashio::fs.file_exists '/config/settings.js'; then
    mkdir -p /config/nodes \
        || bashio::exit.nok "Failed to create node-red configuration directory"

    # Copy in template files & Create random flow id
    id=$(node -e "console.log((2+Math.random()*4294967295).toString(26));")
    if ! bashio::fs.file_exists '/config/flows.json'; then
        cp /etc/node-red/flows.json /config/ && sed -i "s/%%ID%%/${id}/" "/config/flows.json" || bashio::log.info "Failed cp /etc/node-red/flows.json /config/ template"
    fi
    cp /etc/node-red/settings.js /config/ || bashio::log.info "Failed cp /etc/node-red/settings.js /config/ template"
fi

# Pass in port & SSL settings
dashboard_ui_base="/opt/node_modules/node-red-dashboard/nodes/ui_base.html"
port=$(bashio::addon.port 80)
if bashio::fs.file_exists "${dashboard_ui_base}"; then
    sed -i "s/%%PORT%%/${port:-80}/" "${dashboard_ui_base}"
    if ! bashio::var.has_value "${port}"; then
        bashio::log.warning
        bashio::log.warning "Direct access mode is disabled, Node-RED Dashboard"
        bashio::log.warning "will not work!"
        bashio::log.warning
        bashio::log.warning "Please assign a port in the Network section of this"
        bashio::log.warning "add-on configuration."
        bashio::log.warning
    fi

    if bashio::config.true 'ssl'; then
        sed -i "s/%%SSL%%/true/" "${dashboard_ui_base}"
    else
        sed -i "s/%%SSL%%/false/" "${dashboard_ui_base}"
    fi
else
    bashio::log.warning "node-red-dashboard UI file not found; skipping dashboard patch setup"
fi

# Ensures conflicting Node-RED packages are absent
cd /config || bashio::exit.nok "Could not change directory to Node-RED"
if bashio::fs.file_exists "/config/package.json"; then
    npm uninstall \
        node-red-contrib-home-assistant \
        node-red-contrib-home-assistant-llat \
        node-red-contrib-home-assistant-ws \
            || bashio::log.warning "Failed un-installing conflicting packages"
fi
